# trainerright.path
[Path]
PathExists=/var/lib/docker/volumes/trainerIBOS_VOLUME

[Install]
WantedBy=multi-user.target
