# trainerright.service
[Service] 
ExecStart=chmod -R o+rwx /var/lib/docker/volumes

[Install]
WantedBy=multi-user.target
