# setup/units/Lab1.service
[Unit]
Description= CheckLab1

[Service]
ExecStart=python3 {0}

[Install]
WantedBy=multi-user.target