# setup/units/Lab4No.service
[Unit]
Description= CheckLab4

[Service]
ExecStart=python3 {0}

[Install]
WantedBy=multi-user.target