# setup/units/Lab5.service
[Unit]
Description= CheckLab5

[Service]
ExecStart=python3 {0}

[Install]
WantedBy=multi-user.target