# images/Dockerfiles/5/ProgrammOne.service
[Unit]
Description=ProgrammOne for Journald Lab. Good Luck you)

[Service]
ExecStart=/mnt/ExecProgrammOne

[Install]
WantedBy=multi-user.target