# images/Dockerfiles/5/ProgrammTwo.service
[Unit]
Description=ProgrammTwo for Journald Lab. Nice work!

[Service]
ExecStart=/mnt/ExecProgrammTwo

[Install]
WantedBy=multi-user.target